# com.tradingbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 服务标识符 -->
    <key>Label</key>
    <string>com.tradingbot</string>
    
    <!-- 程序路径和参数 -->
    <key>ProgramArguments</key>
    <array>
        <!-- Python解释器路径 - 需要修改为实际路径 -->
        <string>/path/to/your/project/.venv/bin/python</string>
        <!-- 交易主程序路径 - 需要修改为实际路径 -->
        <string>/path/to/your/project/src/trading_main.py</string>
    </array>
    
    <!-- 工作目录 -->
    <key>WorkingDirectory</key>
    <string>/path/to/your/project</string>
    
    <!-- 加载时启动 -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- 保持运行（崩溃后自动重启） -->
    <key>KeepAlive</key>
    <true/>
    
    <!-- 标准输出日志 - 需要修改为实际路径 -->
    <key>StandardOutPath</key>
    <string>/path/to/your/project/logs/trading.stdout.log</string>
    
    <!-- 标准错误日志 - 需要修改为实际路径 -->
    <key>StandardErrorPath</key>
    <string>/path/to/your/project/logs/trading.stderr.log</string>
    
    <!-- 环境变量 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
